# docs/qbittorrent-protonvpn/protonvpn-portforward.service
[Unit]
Description=ProtonVPN NAT-PMP Port Forwarding for qBittorrent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/root/unifi-portforward.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=protonvpn-portforward

[Install]
WantedBy=multi-user.target
